# system-backup.service
[Unit]
Description=Daily System Backup

[Service]
ExecStart=%h/scripts/system-backup.sh
Nice=10
